# supervisord_redmine.conf
[program:redmine]
command = /home/redmine/redmine/startup.sh
